# .config/systemd/user/variety.service
[Unit]
Description=Variety wallpaper rotation tool

[Service]
ExecStart=/usr/bin/variety

[Install]
WantedBy=default.target
